# system_linux/system/jiotv_go.service.d/user.conf
# Override service user
[Service]
User=jiotv_go
Group=jiotv_go
